# mewgrepd.service
[Unit]
Description=mewgrepd

[Service]
ExecStart=/home/masm/bin/mewgrepd.py
Restart=always

[Install]
WantedBy=default.target
